# woof-code/rootfs-petbuilds/pcmanfm/etc/xdg/autostart/pcmanfm.desktop
[Desktop Entry]
Version=1.0
Name=pcmanfm
Comment=pcmanfm
Exec=pcmanfm --desktop
Terminal=false
Type=Application