# woof-code/rootfs-petbuilds/busybox/etc/xdg/autostart/acpid.desktop
[Desktop Entry]
Version=1.0
Name=acpid
Comment=ACPI daemon
Exec=acpid
Terminal=false
Type=Application