# woof-code/rootfs-skeleton/etc/udev/rules.d/91-xinput-apply.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/bin/sh -c 'sleep 1; xinput-apply'"
