# woof-code/rootfs-skeleton/etc/udev/rules.d/92-xkb-apply.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/bin/sh -c 'sleep 1; xkb-apply'"
