# woof-code/rootfs-petbuilds/freememapplet/root/.config/autostart/freememapplet.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
NoDisplay=true
Name=freememapplet_tray
Exec=freememapplet_tray
