# woof-code/rootfs-petbuilds/dwl-kiosk/etc/xdg/autostart/dwl-watermark.desktop
[Desktop Entry]
Version=1.0
Name=dwl-watermark
Comment=Development build watermark
Exec=dwl-watermark
Terminal=false
Type=Application
