# woof-code/rootfs-petbuilds/yambar-dwl/etc/xdg/autostart/yambar.desktop
[Desktop Entry]
Version=1.0
Name=yambar
Comment=yambar
Exec=yambar
Terminal=false
Type=Application
