# woof-code/rootfs-packages/acpid_busybox/etc/xdg/autostart/acpid.desktop
[Desktop Entry]
Version=1.0
Name=acpid
Comment=ACPI daemon
Exec=startacpid
Terminal=false
Type=Application