# woof-code/rootfs-petbuilds/spot-pkexec/etc/xdg/autostart/pkexecd.desktop
[Desktop Entry]
Version=1.0
Name=pkexecd
Comment=pkexec daemon
Exec=pkexecd
Terminal=false
Type=Application
