# woof-code/rootfs-petbuilds/yambar-puppy/etc/xdg/autostart/yambar.desktop
[Desktop Entry]
Version=1.0
Name=yambar
Comment=yambar
Exec=startyambar
Terminal=false
Type=Application
