# woof-code/rootfs-skeleton/etc/xdg/autostart/mscw.desktop
[Desktop Entry]
Type=Application
Name=Multiple Sound Card Wizard
NoDisplay=true
Exec=sh -c 'sleep 5; mscw start'
