# woof-code/rootfs-skeleton/etc/xdg/autostart/xinput-apply.desktop
[Desktop Entry]
Type=Application
Name=Apply libinput settings
NoDisplay=true
Exec=xinput-apply
