# woof-code/rootfs-skeleton/etc/xdg/autostart/xkb-apply.desktop
[Desktop Entry]
Type=Application
Name=Apply keyboard settings
NoDisplay=true
Exec=xkb-apply
