# woof-code/rootfs-skeleton/etc/xdg/autostart/xorg-mousekeys.desktop
[Desktop Entry]
Type=Application
Name=Enable Xorg Mouse Keys
NoDisplay=true
Exec=sh -c 'setxkbmap -option keypad:pointerkeys'

