# woof-code/woof2-scripts_and_files/packages-templates/hal/etc/udev/rules.d/90-hal.rules
# BK trying to emasculate HAL, commented out...
## pass all events to the HAL daemon
# RUN+="socket:/org/freedesktop/hal/udev_event"
