# woof-code/rootfs-skeleton/etc/profile.d/flatpak-bindir.sh
[ -e /var/lib/flatpak/exports/bin ] && export PATH="$PATH:/var/lib/flatpak/exports/bin"