# woof-code/rootfs-skeleton/etc/profile.d/glib-memory-alloc
#!/bin/sh

export G_SLICE="always_malloc"
