# woof-code/rootfs-skeleton/etc/profile.d/gtkdialog
if ldd $(which gtkdialog) | grep -Fq 'gtk-3' ; then #gtk3
	GTKDIALOG_BUILD=GTK3
elif ldd $(which gtkdialog) | grep -Fq 'gtk-X11-2' ; then #gtk2
	GTKDIALOG_BUILD=GTK2
fi
export GTKDIALOG_BUILD
