# woof-code/rootfs-skeleton/etc/profile.d/pup_gtk
#Qt4 apps need this variable (also see Style=GTK+ in /root/.config/Trolltech.conf)...
#Qt4 note also, /root/.gtkrc-2.0 needs this line: gtk-theme-name="name of gtk theme"
export GTK2_RC_FILES=/etc/gtk-2.0/gtkrc
