# woof-code/rootfs-skeleton/etc/profile.d/qt-libs
#!/bin/sh

for libp in /usr/lib /usr/lib32 /usr/lib/i386-linux-gnu /usr/local/lib /usr/local/lib32 /usr/local/lib/i386-linux-gnu 
do 

	if [ -e $libp/qt ];then
		
		[ -d $libp/qt/bin ] && PATH="${PATH}:$libp/qt/bin"
		[ -d $libp/qt/lib ] && LD_LIBRARY_PATH="$LD_LIBRARY_PATH:$libp/qt/lib"
		
		if [ "$QTDIR" == "" ]; then
		 QTDIR=$libp/qt
		else
		 QTDIR=$QTDIR:$libp/qt
		fi
		
	fi

done

for libp in /usr/lib64 /usr/lib/x86_64 /usr/lib/x64 /usr/lib/x86_64-linux-gnu /usr/local/lib64 /usr/local/lib/x86_64 /usr/local/lib/x64 /usr/lib/local/x86_64-linux-gnu
do 

	if [ -e $libp/qt ];then
		
		[ -d $libp/qt/bin ] && PATH="${PATH}:$libp/qt/bin"
		[ -d $libp/qt/lib ] && LD_LIBRARY_PATH="$LD_LIBRARY_PATH:$libp/qt/lib"
		
		if [ "$QTDIR" == "" ]; then
		 QTDIR=$libp/qt
		else
		 QTDIR=$QTDIR:$libp/qt
		fi
		
	fi

done

if [ -d "/opt/qt" ];then
	
	[ -d /opt/qt/lib ] && LD_LIBRARY_PATH="$LD_LIBRARY_PATH:/opt/qt/lib"
	[ -d /opt/qt/bin ] && PATH="$PATH:/opt/qt/bin" #if devx sfs.
	
	if [ "$QTDIR" == "" ]; then
	 QTDIR="/opt/qt"
	else
	 QTDIR="$QTDIR:/opt/qt"
	fi
	
	if [ -d /opt/qt/include ]; then	 
		if [ "$CPLUS_INCLUDE_PATH" == "" ]; then
		 export CPLUS_INCLUDE_PATH="/opt/qt/include" #devx
		else
		  CPLUS_INCLUDE_PATH="$CPLUS_INCLUDE_PATH:/opt/qt/include"
		fi
	fi
		
fi

export QTDIR


for qtver in 2 3 4 5 6 7 8 9 10
do

    for libp in /usr/lib /usr/lib32 /usr/lib/i386-linux-gnu /usr/local/lib /usr/local/lib32 /usr/local/lib/i386-linux-gnu
	do 
	
	 if [ -e $libp/qt${qtver} ];then
		
		[ -d $libp/qt${qtver}/bin ] && PATH="${PATH}:$libp/qt${qtver}/bin"
		[ -d $libp/qt${qtver}/lib ] && LD_LIBRARY_PATH="$LD_LIBRARY_PATH:$libp/qt${qtver}/lib"
	
		if [ "$(eval "echo \$QT${qtver}DIR")" == "" ]; then
		 eval "QT${qtver}DIR=\"$libp/qt${qtver}\""
		else
		 eval "QT${qtver}DIR=\"\$QT${qtver}DIR:$libp/qt${qtver}\""
		fi
	
	 fi

    done
    
    
    for libp in /usr/lib64 /usr/lib/x86_64 /usr/lib/x64 /usr/lib/x86_64-linux-gnu /usr/local/lib64 /usr/local/lib/x86_64 /usr/local/lib/x64 /usr/lib/local/x86_64-linux-gnu
	do 
	
	 if [ -e $libp/qt${qtver} ];then
		
		[ -d $libp/qt${qtver}/bin ] && PATH="${PATH}:$libp/qt${qtver}/bin"
		[ -d $libp/qt${qtver}/lib ] && LD_LIBRARY_PATH="$LD_LIBRARY_PATH:$libp/qt${qtver}/lib"
	
		if [ "$(eval "echo \$QT${qtver}DIR")" == "" ]; then
		 eval "QT${qtver}DIR=\"$libp/qt${qtver}\""
		else
		 eval "QT${qtver}DIR=\"\$QT${qtver}DIR:$libp/qt${qtver}\""
		fi
	
	 fi

    done
    

	if [ -d "/opt/qt${qtver}" ];then
	
		[ -d /opt/qt${qtver}/bin ] && PATH="$PATH:/opt/qt${qtver}/bin" #if devx sfs.
	    [ -d /opt/qt${qtver}/lib ] && LD_LIBRARY_PATH="$LD_LIBRARY_PATH:/opt/qt${qtver}/lib"
		
		if [ "$(eval "echo \$QT${qtver}DIR")" == "" ]; then
		 eval "QT${qtver}DIR=\"/opt/qt${qtver}\""
		else
		 eval "QT${qtver}DIR=\"\$QT${qtver}DIR:/opt/qt${qtver}\""
		fi
		
		if [ -d /opt/qt${qtver}/include ]; then		 
			if [ "$CPLUS_INCLUDE_PATH" == "" ]; then
			  CPLUS_INCLUDE_PATH="/opt/qt${qtver}/include" #devx
			else
			  CPLUS_INCLUDE_PATH="$CPLUS_INCLUDE_PATH:/opt/qt${qtver}/include"
			fi
		fi
			
	fi
	
	eval "export QT${qtver}DIR"
 
done

export CPLUS_INCLUDE_PATH

export QT_XFT=true
