# woof-code/rootfs-skeleton/etc/inittab
::sysinit:/etc/rc.d/rc.sysinit
tty1::respawn:plogin
tty2::respawn:getty 38400 tty2
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting

# dialog misbehaves if rc.shutdown is called by init (busybox reboot/poweroff)
#::shutdown:/etc/rc.d/rc.shutdown

::shutdown:/etc/rc.d/rc.cleanup
