# woof-code/rootfs-skeleton/etc/inittab
::sysinit:/etc/rc.d/rc.sysinit
tty1::respawn:plogin
tty2::respawn:getty 38400 tty2
tty3::respawn:getty 38400 tty3
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/etc/rc.d/rc.shutdown
::shutdown:/etc/rc.d/rc.cleanup
