# woof-code/rootfs-packages/simple_network_setup/etc/udev/rules.d/51-simple_network_setup.rules
# Accumulate names of modules to be loaded, for use by rc.network -- if pup_event_backend_modprobe installed.
#Note that for ssb (and probably brcm) subsystem, %s{modalias} = null, so is unusable.
#Note that ENV{MODALIAS}=="?*" may kill the rule in later kernels (e.g., for VoidPup).

SUBSYSTEM=="pci|usb|ssb|bcma|pcmcia", ACTION=="add", \
  RUN+="/usr/local/simple_network_setup/build_udevmodulelist"
