# woof-code/rootfs-packages/simple_network_setup/etc/udev/rules.d/51-simple_network_setup.rules
# Accumulate names of modules to be loaded, for use by rc.network -- if pup_event_backend_modprobe installed.

PROGRAM=="/usr/bin/test -x /sbin/pup_event_backend_modprobe", \
  PROGRAM=="/usr/bin/test -s /etc/simple_network_setup/connections", \
  SUBSYSTEM=="?*", SUBSYSTEM!="acpi|input|platform", \
  ENV{MODALIAS}=="?*", ACTION=="add", \
  RUN+="/usr/local/simple_network_setup/build_udevmodulelist $env{MODALIAS}"
