# woof-code/rootfs-packages/simple_network_setup/etc/udev/rules.d/51-simple_network_setup.rules
# Accumulate list of modules to be loaded, for use by rc.network.

ENV{MODALIAS}=="?*", ACTION=="add", SUBSYSTEM=="?*", \
  PROGRAM="/usr/bin/test -s /etc/simple_network_setup/connections", \
  RUN+="/usr/local/simple_network_setup/build_udevmodulelist"
