# woof-code/rootfs-petbuilds/foot-puppy/etc/xdg/autostart/foot-server.desktop
[Desktop Entry]
Type=Application
Name=Foot server
NoDisplay=true
Exec=foot -s
