# woof-code/rootfs-petbuilds/labwc/etc/xdg/autostart/altf1d.desktop
[Desktop Entry]
Version=1.0
Name=altf1d
Comment=Alt+F1 daemon
Exec=altf1d
Terminal=false
Type=Application